# config/macos/com.librarybabel.api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.librarybabel.api</string>
    
    <key>Program</key>
    <string>/opt/homebrew/bin/python3</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/weixiangzhang/Local Dev/LibraryOfBabel/src/api/production_api.py</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/weixiangzhang/Local Dev/LibraryOfBabel</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/weixiangzhang/Local Dev/LibraryOfBabel/logs/api.out.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/weixiangzhang/Local Dev/LibraryOfBabel/logs/api.err.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PYTHONPATH</key>
        <string>/Users/weixiangzhang/Local Dev/LibraryOfBabel/src</string>
        <key>API_KEY</key>
        <string>babel_secure_8a52a0ad3a1fe3bf3ade37d04deef0054d8f58035a0e9d4760a9a08548d8cebf</string>
        <key>OLLAMA_MODEL</key>
        <string>llama3.2:3b</string>
    </dict>
    
    <key>ThrottleInterval</key>
    <integer>5</integer>
    
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>